# daemon/foxd-daemon.service
[Unit]
Description=Fox Daemon - LAN Monitoring Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root

# Path to the foxd-daemon binary
ExecStart=/usr/local/bin/foxd-daemon

# Working directory
WorkingDirectory=/etc/foxd

# Environment variables
Environment="FOXD_CONFIG=/etc/foxd/config.toml"
Environment="RUST_LOG=foxd_daemon=info"

# Restart policy
Restart=always
RestartSec=10

# Security settings
# Note: Root is required for packet capture, but you can drop capabilities after startup
# with a more sophisticated setup
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN
NoNewPrivileges=true

# Resource limits
LimitNOFILE=65536

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=foxd-daemon

[Install]
WantedBy=multi-user.target
